# deploy/systemd/aot-cluster.target
[Unit]
Description=AOT Development Cluster (all services)

[Install]
WantedBy=default.target
